# config_files/etc/systemd/system/scibot-bookmarklet.socket
[Unit]
Description=scibot bookmarklet gunicorn socket

[Socket]
ListenStream=/run/scibot-bookmarket/socket

[Install]
WantedBy=sockets.target
