# config_files/etc/systemd/system/scibot-dashboard.socket
[Unit]
Description=scibot dashboard gunicorn socket

[Socket]
ListenStream=/run/scibot-dashboard/socket

[Install]
WantedBy=sockets.target
